# smartcard-diagnostic/daemon/com.canonical.SmartCardMonitor.service
[Unit]
Description=Smart Card Monitor D-Bus Service
After=graphical-session.target

[Service]
Type=dbus
BusName=com.canonical.SmartCardMonitor
ExecStart=/snap/smartcard-diagnostic/current/bin/smartcard-monitor
Restart=on-failure
RestartSec=5

[Install]
WantedBy=default.target
