# smartcard-diagnostic/daemon/com.canonical.SmartCardMonitor.service
[Unit]
Description=Smart Card Monitor D-Bus Service
After=pcscd.service pcscd.socket

[Service]
Type=dbus
BusName=com.canonical.SmartCardMonitor
ExecStart=/snap/smartcard-diagnostic/current/bin/smartcard-monitor
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
